# deploy/ipcap-agent-listen.service
[Unit]
Description=ipcap Noise listener (serves the spool to the collector)
Documentation=https://github.com/MadrHacks/ipcap
After=network-online.target ipcap-agent.service
Wants=network-online.target

[Service]
# Read-only Noise responder, crash-isolated from the capture writer. Authenticates
# the collector by static key; --peer is the collector's base64 public key
# (repeatable for multiple collectors). The agent's own private key is generated
# once with `ipcap keygen --out /etc/ipcap/agent.key` (mode 0600); its public key
# goes into the collector's vulnbox.yml as noise_pubkey.
ExecStart=/usr/local/bin/ipcap agent listen \
    --spool-dir /var/lib/ipcap/spool \
    --src-id 1 \
    --listen :7878 \
    --key /etc/ipcap/agent.key \
    --peer REPLACE_WITH_COLLECTOR_PUBLIC_KEY
Restart=always
RestartSec=1
StartLimitIntervalSec=0
# No privileges needed: it only reads the spool and a TCP socket.
DynamicUser=no
User=root
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
