# deploy/ipcap-agent.service
[Unit]
Description=ipcap durable capture agent
Documentation=https://github.com/MadrHacks/ipcap
After=network-online.target
Wants=network-online.target

[Service]
# Persistent capture: the zero-loss anchor, decoupled from any collector.
# %i is the capture interface (e.g. ipcap-agent@eth0).
#
# IMPORTANT (anti-amplification): --ssh-port MUST equal this host's real sshd
# port (the one the collector drains over), and/or pass --mgmt <collector-host-IP>
# (repeatable, accepts IPs or CIDRs). If the drain's SSH traffic is not excluded
# it gets captured, spooled, and re-served — recaptured again — amplifying without
# bound. Templated by the ansible playbook from vulnbox.yml's ssh_port (M3).
ExecStart=/usr/local/bin/ipcap agent capture \
    --iface %i \
    --spool-dir /var/lib/ipcap/spool \
    --src-id 1 \
    --ssh-port 22 \
    --mgmt 0.0.0.0/32 \
    --retention-bytes 51539607552 \
    --ring-mib 256
# Capture needs raw sockets; nothing else.
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
# Survive sub-second restarts; the large ring covers the gap. Config reload must
# NOT restart capture (only the serve/collector side reloads).
Restart=always
RestartSec=1
StartLimitIntervalSec=0
StateDirectory=ipcap
# notify + watchdog: the capture loop pings the watchdog, so a silently-wedged
# capture (process alive, loop stalled) is detected and restarted within 30s.
Type=notify
WatchdogSec=30

[Install]
WantedBy=multi-user.target
